# debian/mmta.cron.d
7,37 * * * *   root  /usr/share/mmta/mmta-send-allusers
